# files/etc/init.d/handle-keys
#!/bin/sh /etc/rc.common

START=99

start() {
    /usr/bin/handle-keys.sh &
}

stop() {
    killall handle-keys.sh
}
